# services/dbus-ibr-system/run
#!/bin/sh
exec 2>&1
exec /data/ibr-venus-service/dbus-ibr-system/dbus-ibr-system.sh
