# services/dbus-ibr-system/run
#!/bin/sh
exec 2>&1
exec /data/ibr-venus-services/dbus-ibr-system/dbus-ibr-system.sh
